# 部分代码/wechat/rc.local
#!bin/bash -e
/data/apps/wechat/get_moment_loop >/dev/null 2>&1 & 
